# deploy/launchd/io.tabura.llm.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>io.tabura.llm</string>
  <key>ProgramArguments</key>
  <array>
    <string>@@LLM_SETUP_SCRIPT@@</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>TABURA_LLM_PRESET</key>
    <string>fast-qwen9b</string>
    <key>TABURA_LLM_MODEL_DIR</key>
    <string>@@LLM_MODEL_DIR@@</string>
    <key>TABURA_LLM_MODEL_FILE</key>
    <string>Qwen3.5-9B-Q4_K_M.gguf</string>
    <key>TABURA_LLM_MODEL_URL</key>
    <string>https://huggingface.co/lmstudio-community/Qwen3.5-9B-GGUF/resolve/main/Qwen3.5-9B-Q4_K_M.gguf?download=true</string>
    <key>TABURA_LLM_CTX</key>
    <string>32768</string>
    <key>TABURA_LLM_THREADS</key>
    <string>8</string>
    <key>LLAMA_SERVER_BIN</key>
    <string>@@LLAMA_SERVER_BIN@@</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/tabura-llm.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/tabura-llm.log</string>
</dict>
</plist>
