# deploy/launchd/io.tabura.vllm.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>io.tabura.vllm</string>
  <key>ProgramArguments</key>
  <array>
    <string>@@VLLM_SETUP_SCRIPT@@</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>TABURA_VLLM_ROOT</key>
    <string>@@VLLM_ROOT@@</string>
    <key>TABURA_VLLM_HOST</key>
    <string>127.0.0.1</string>
    <key>TABURA_VLLM_PORT</key>
    <string>8426</string>
    <key>TABURA_VLLM_MODEL</key>
    <string>QuantTrio/Qwen3.5-9B-AWQ</string>
    <key>TABURA_VLLM_SERVED_MODEL_NAME</key>
    <string>tabura-qwen-9b</string>
    <key>TABURA_VLLM_MAX_MODEL_LEN</key>
    <string>4096</string>
    <key>TABURA_VLLM_GPU_MEMORY_UTILIZATION</key>
    <string>0.75</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/tabura-vllm.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/tabura-vllm.log</string>
</dict>
</plist>
