# deploy/systemd/user/tabura-stt.service
[Unit]
Description=Tabura STT Sidecar (Voxtral preferred, voxtype fallback)
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/code/assi/tabula
Environment=HOME=%h
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=TABURA_STT_PROVIDER=auto
Environment=TABURA_STT_LANGUAGE=de,en
Environment=TABURA_STT_MODEL=large-v3-turbo
Environment=TABURA_VOXTRAL_MODEL=mistralai/Voxtral-Mini-3B-2507
Environment=TABURA_VOXTRAL_SERVED_MODEL=whisper-1
ExecStart=%h/code/assi/tabula/scripts/setup-stt-sidecar.sh
Restart=on-failure
RestartSec=5
TimeoutStopSec=15

[Install]
WantedBy=default.target
