# deploy/systemd/user/tabura-vllm.service
[Unit]
Description=Tabura Local vLLM Runtime
After=network.target

[Service]
Type=simple
Environment=TABURA_VLLM_ROOT=%h/.local/share/tabura-vllm
Environment=TABURA_VLLM_HOST=127.0.0.1
Environment=TABURA_VLLM_PORT=8426
Environment=TABURA_VLLM_MODEL=QuantTrio/Qwen3.5-9B-AWQ
Environment=TABURA_VLLM_SERVED_MODEL_NAME=tabura-qwen-9b
Environment=TABURA_VLLM_MAX_MODEL_LEN=4096
Environment=TABURA_VLLM_GPU_MEMORY_UTILIZATION=0.75
ExecStart=@@REPO_ROOT@@/scripts/setup-local-vllm.sh
Restart=on-failure
RestartSec=5
TimeoutStopSec=30

[Install]
WantedBy=default.target
